# scripts/services/raspdash-rfcomm.service
[Unit]
Description=RaspDash ELM327 Bluetooth RFCOMM binding
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=-/usr/bin/rfcomm release 0
ExecStart=/usr/bin/rfcomm bind 0 __ELM327_MAC__ 1
ExecStop=/usr/bin/rfcomm release 0

[Install]
WantedBy=multi-user.target
